# backend/app/reference_data/attack.json
[
  {"id": "T1566", "name": "Phishing", "tactic": "Initial Access", "description": "Adversaries send phishing messages to gain access to victim systems"},
  {"id": "T1566.001", "name": "Spearphishing Attachment", "tactic": "Initial Access", "description": "Adversaries send spearphishing emails with a malicious attachment"},
  {"id": "T1566.002", "name": "Spearphishing Link", "tactic": "Initial Access", "description": "Adversaries send spearphishing emails with a malicious link"},
  {"id": "T1190", "name": "Exploit Public-Facing Application", "tactic": "Initial Access", "description": "Adversaries exploit a weakness in an internet-facing application"},
  {"id": "T1133", "name": "External Remote Services", "tactic": "Initial Access", "description": "Adversaries leverage external-facing remote services to access a network"},
  {"id": "T1195", "name": "Supply Chain Compromise", "tactic": "Initial Access", "description": "Adversaries manipulate products or delivery mechanisms prior to receipt by the end consumer"},
  {"id": "T1078", "name": "Valid Accounts", "tactic": "Defense Evasion", "description": "Adversaries obtain and abuse credentials of existing accounts"},
  {"id": "T1078.004", "name": "Cloud Accounts", "tactic": "Defense Evasion", "description": "Adversaries obtain and abuse credentials of cloud accounts"},
  {"id": "T1059", "name": "Command and Scripting Interpreter", "tactic": "Execution", "description": "Adversaries abuse command and script interpreters to execute commands"},
  {"id": "T1059.001", "name": "PowerShell", "tactic": "Execution", "description": "Adversaries abuse PowerShell for execution and scripting"},
  {"id": "T1203", "name": "Exploitation for Client Execution", "tactic": "Execution", "description": "Adversaries exploit software vulnerabilities in client applications"},
  {"id": "T1053", "name": "Scheduled Task/Job", "tactic": "Persistence", "description": "Adversaries abuse task scheduling to establish persistence"},
  {"id": "T1136", "name": "Create Account", "tactic": "Persistence", "description": "Adversaries create accounts to maintain access to victim systems"},
  {"id": "T1547", "name": "Boot or Logon Autostart Execution", "tactic": "Persistence", "description": "Adversaries configure system settings to automatically execute a program during boot or logon"},
  {"id": "T1548", "name": "Abuse Elevation Control Mechanism", "tactic": "Privilege Escalation", "description": "Adversaries circumvent mechanisms designed to control elevated privileges"},
  {"id": "T1068", "name": "Exploitation for Privilege Escalation", "tactic": "Privilege Escalation", "description": "Adversaries exploit software vulnerabilities to elevate privileges"},
  {"id": "T1134", "name": "Access Token Manipulation", "tactic": "Privilege Escalation", "description": "Adversaries modify access tokens to operate under a different security context"},
  {"id": "T1562", "name": "Impair Defenses", "tactic": "Defense Evasion", "description": "Adversaries maliciously modify components of a victim environment to hinder defenses"},
  {"id": "T1070", "name": "Indicator Removal", "tactic": "Defense Evasion", "description": "Adversaries delete or modify artifacts generated on a host system"},
  {"id": "T1027", "name": "Obfuscated Files or Information", "tactic": "Defense Evasion", "description": "Adversaries attempt to make payloads difficult to analyze by encoding or encrypting them"},
  {"id": "T1110", "name": "Brute Force", "tactic": "Credential Access", "description": "Adversaries use brute force techniques to attempt access when passwords are unknown"},
  {"id": "T1003", "name": "OS Credential Dumping", "tactic": "Credential Access", "description": "Adversaries attempt to dump credentials to obtain account login material"},
  {"id": "T1528", "name": "Steal Application Access Token", "tactic": "Credential Access", "description": "Adversaries steal application access tokens as a means of acquiring credentials"},
  {"id": "T1552", "name": "Unsecured Credentials", "tactic": "Credential Access", "description": "Adversaries search for common password storage locations to obtain user credentials"},
  {"id": "T1087", "name": "Account Discovery", "tactic": "Discovery", "description": "Adversaries attempt to get a listing of valid accounts on a system or domain"},
  {"id": "T1046", "name": "Network Service Discovery", "tactic": "Discovery", "description": "Adversaries attempt to get a listing of services running on remote hosts"},
  {"id": "T1057", "name": "Process Discovery", "tactic": "Discovery", "description": "Adversaries attempt to get information about running processes"},
  {"id": "T1021", "name": "Remote Services", "tactic": "Lateral Movement", "description": "Adversaries use valid accounts to log into remote services for lateral movement"},
  {"id": "T1550", "name": "Use Alternate Authentication Material", "tactic": "Lateral Movement", "description": "Adversaries use alternate authentication material such as hashes or tickets"},
  {"id": "T1560", "name": "Archive Collected Data", "tactic": "Collection", "description": "Adversaries compress and encrypt data prior to exfiltration"},
  {"id": "T1005", "name": "Data from Local System", "tactic": "Collection", "description": "Adversaries search local system sources for data of interest prior to exfiltration"},
  {"id": "T1039", "name": "Data from Network Shared Drive", "tactic": "Collection", "description": "Adversaries search network shares for files of interest"},
  {"id": "T1041", "name": "Exfiltration Over C2 Channel", "tactic": "Exfiltration", "description": "Adversaries steal data by exfiltrating it over an existing command and control channel"},
  {"id": "T1567", "name": "Exfiltration Over Web Service", "tactic": "Exfiltration", "description": "Adversaries use legitimate external web services to exfiltrate data"},
  {"id": "T1486", "name": "Data Encrypted for Impact", "tactic": "Impact", "description": "Adversaries encrypt data on target systems to interrupt availability"},
  {"id": "T1489", "name": "Service Stop", "tactic": "Impact", "description": "Adversaries stop or disable services on a system to render those services unavailable"},
  {"id": "T1498", "name": "Network Denial of Service", "tactic": "Impact", "description": "Adversaries perform network denial of service to degrade or block availability"},
  {"id": "T1071", "name": "Application Layer Protocol", "tactic": "Command and Control", "description": "Adversaries communicate using OSI application layer protocols to avoid detection"},
  {"id": "T1105", "name": "Ingress Tool Transfer", "tactic": "Command and Control", "description": "Adversaries transfer tools or files from an external system into a compromised environment"},
  {"id": "T1219", "name": "Remote Access Software", "tactic": "Command and Control", "description": "Adversaries use legitimate remote access tools to establish an interactive C2 channel"}
]
